# .etc/cron.daily/zzz_clean_var_log
#!/bin/sh

find /var/log \( -name '*.gz' -o -name '*.1' \) -exec rm {} \;
